# meta-renesas/meta-rz-common/recipes-graphics/dummy-input/dummy-input/dummy-input.service
[Unit]
Description=dummy-input

[Service]
Type=simple
ExecStart=@BINDIR@/dummy-input

[Install]
WantedBy=multi-user.target
